# scripts/gptbot-api.service
[Unit]
Description=GPT Bot API Server
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=gptbot
Group=gptbot
WorkingDirectory=/opt/gptbot
Environment="PATH=/opt/gptbot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/gptbot/venv/bin/python -m app.api.run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

